# OVMS/ovms.service
[Unit]
Description=OVMS Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rune
WorkingDirectory=/home/rune/claude_workspace/OVMS/webapp
ExecStart=/usr/bin/python3 serve.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
